# templates/systemd/sonarr.service.j2
[Unit]
Description=Sonarr Daemon.
After=network.target

[Service]
Type=simple
ExecStart={{ sonarr__cfg._app ~ "/Sonarr" }} -nobrowser -data={{ sonarr__cfg._var_dir }}
TimeoutStopSec=20
Restart=always
User={{ sonarr__srv._user }}
Group={{ sonarr__srv._group }}
UMask=0002
KillMode=process

[Install]
WantedBy=multi-user.target